# systemd/vpn-update.service
[Unit]
Description=Run VPN update

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vpn update
