# systemd/vpn-reset.service
[Unit]
Description=Run VPN init-reset at boot

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vpn init-reset

[Install]
WantedBy=multi-user.target
