# systemd/vpn-connect.service
[Unit]
Description=Run VPN init-connect at boot
After=network-online.target vpn-killswitch.service vpn-reset.service
Wants=network-online.target vpn-killswitch.service vpn-reset.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vpn init-connect

[Install]
WantedBy=multi-user.target
